# systemd/nanobot-failure-notify.service
[Unit]
Description=NanoBot failure Telegram notification

[Service]
Type=oneshot
ExecStart=/bin/bash %h/claude-tg-tmux/notify-failure.sh
